# roles/update/files/stk-check.timer
[Unit]
Description=Daily check if system components are up to date

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
