# roles/xdg_open/templates/ssh_remote_xdg_open@.service.j2
[Unit]
Description=Local "open URL" service instance

[Service]
ExecStart=/bin/bash -c 'IFS= read -r -d "" url && {{ xdg_open_open_cmd }} "$url"'
StandardInput=socket
