# roles/xdg_open/templates/ssh_remote_xdg_open@.service.j2
[Unit]
Description=Local "open URL" service instance

[Service]
ExecStart=/bin/sh -lc 'xargs -0 -n1 xdg-open'
StandardInput=socket
